# ssh/rootfs/etc/services.d/stdin/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Home Assistant Community Add-on: SSH & Web Terminal
# Starts the STDIN service
# ==============================================================================

bashio::log.info 'Starting the Home Assistant STDIN service...'

# shellcheck disable=SC2162
while read cmd; do
  cmd="${cmd%\"}"
  cmd="${cmd#\"}"

  bashio::log.info 'Received external command to execute:'
  bashio::log.info "> ${cmd}"

  if bash --login -c "${cmd}"; then
    bashio::log.info "Command completed successfully."
  else
    bashio::log.error "Command exited with a failure status."
  fi
done < /proc/1/fd/0
