# ssh/rootfs/etc/s6-overlay/s6-rc.d/init-user/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Community Add-on: Advanced SSH & Web Terminal
# Executes configured customizations & persists user settings
# ==============================================================================
readonly -a DIRECTORIES=(addon_configs addons backup homeassistant media share ssl)
readonly BASH_HISTORY_FILE=/root/.bash_history
readonly BASH_HISTORY_PERSISTENT_FILE=/data/.bash_history
readonly GIT_CONFIG=/data/.gitconfig
readonly HOME_ASSISTANT_PROFILE_D_FILE=/etc/profile.d/homeassistant.sh
readonly SSH_USER_PATH=/data/.ssh
readonly VSCODE_SERVER_PATH=/root/.vscode-server
readonly VSCODE_SERVER_PERSISTENT_PATH=/data/.vscode-server
readonly ZSH_HISTORY_FILE=/root/.zsh_history
readonly ZSH_HISTORY_PERSISTENT_FILE=/data/.zsh_history

# Links some common directories to the user's home folder for convenience
for dir in "${DIRECTORIES[@]}"; do
    ln -s "/${dir}" "${HOME}/${dir}" \
        || bashio::log.warning "Failed linking common directory: ${dir}"
done

# Some links to old locations, to not mess with the user's muscle memory
ln -s "/homeassistant" "/config" \
    || bashio::log.warning "Failed linking common directory: /config"
ln -s "/homeassistant" "${HOME}/config" \
    || bashio::log.warning "Failed linking common directory: ${HOME}/config"

# Sets up ZSH or Bash shell history
if bashio::config.true "zsh"; then
    touch "${ZSH_HISTORY_PERSISTENT_FILE}" \
        || bashio::exit.nok 'Failed creating a persistent ZSH history file'

    chmod 600 "${ZSH_HISTORY_PERSISTENT_FILE}" \
        || bashio::exit.nok \
            'Failed setting the correct permissions to the ZSH history file'

    ln -s -f "${ZSH_HISTORY_PERSISTENT_FILE}" "${ZSH_HISTORY_FILE}" \
        || bashio::exit.nok 'Failed linking the persistent ZSH history file'
else
    touch "${BASH_HISTORY_PERSISTENT_FILE}" \
        || bashio::exit.nok 'Failed creating a persistent Bash history file'

    chmod 600 "${BASH_HISTORY_PERSISTENT_FILE}" \
        || bashio::exit.nok \
            'Failed setting the correct permissions to the Bash history file'

    ln -s -f "${BASH_HISTORY_PERSISTENT_FILE}" "${BASH_HISTORY_FILE}" \
        || bashio::exit.nok 'Failed linking the persistent Bash history file'
fi

# Set up Bash
if ! bashio::config.true "zsh"; then
    sed -i -r -e 's|^(root:.*)/bin/zsh$|\1/bin/bash|' /etc/passwd*
    sed -i -e 's|/zsh$|/bash|' /root/.tmux.conf
fi

echo "export SUPERVISOR_TOKEN=\"${SUPERVISOR_TOKEN}\"" \
    >> "${HOME_ASSISTANT_PROFILE_D_FILE}" \
        || bashio::exit.nok 'Failed to export Supervisor API token'

# Sets up the users .ssh folder to be persistent
if ! bashio::fs.directory_exists "${SSH_USER_PATH}"; then
    mkdir -p "${SSH_USER_PATH}" \
        || bashio::exit.nok 'Failed to create a persistent .ssh folder'

    chmod 700 "${SSH_USER_PATH}" \
        || bashio::exit.nok \
            'Failed setting permissions on persistent .ssh folder'
fi
ln -s "${SSH_USER_PATH}" ~/.ssh

echo "SUPERVISOR_TOKEN=${SUPERVISOR_TOKEN}" > /data/.ssh/environment \
    || bashio::exit.nok 'Failed to set Supervisor API token in env'

# Sets up the user .gitconfig file to be persistent
if ! bashio::fs.file_exists "${GIT_CONFIG}"; then
    touch "${GIT_CONFIG}" \
        || bashio::exit.nok 'Failed to create .gitconfig file'
fi
ln -s "${GIT_CONFIG}" ~/.gitconfig

# Sets up the user .vscode-server folder to be persistent
if ! bashio::fs.directory_exists "${VSCODE_SERVER_PERSISTENT_PATH}"; then
    mkdir -p "${VSCODE_SERVER_PERSISTENT_PATH}" \
        || bashio::exit.nok 'Failed to create a persistent .vscode-server folder'
fi
ln -s "${VSCODE_SERVER_PERSISTENT_PATH}" "${VSCODE_SERVER_PATH}"

# Disable SSH & Web Terminal session sharing if configured
if ! bashio::config.true 'share_sessions'; then
    bashio::log.notice 'Session sharing has been disabled!'
    rm /root/.bash_profile
    rm /root/.zprofile
fi

# Install user configured/requested packages
if bashio::config.has_value 'packages'; then
    apk update \
        || bashio::exit.nok 'Failed updating Alpine packages repository indexes'

    for package in $(bashio::config 'packages'); do
        apk add "$package" \
            || bashio::exit.nok "Failed installing package ${package}"
    done
fi

# Executes user configured/requested commands on startup
if bashio::config.has_value 'init_commands'; then
    while read -r cmd; do
        eval "${cmd}" \
            || bashio::exit.nok "Failed executing init command: ${cmd}"
    done <<< "$(bashio::config 'init_commands')"
fi

# Clone the specified dotfiles repository and execute setup scripts
if bashio::config.has_value 'dotfiles_repo'; then
    dotfiles_repo=$(bashio::config 'dotfiles_repo')
    git clone "${dotfiles_repo}" /root/dotfiles \
        || bashio::exit.nok "Failed to clone dotfiles repository: ${dotfiles_repo}"
    if [ -f /root/dotfiles/install.sh ]; then
        bash /root/dotfiles/install.sh \
            || bashio::exit.nok "Failed to execute dotfiles setup script"
    fi
fi
