# ssh/rootfs/etc/services.d/ttyd/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Community Hass.io Add-ons: SSH & Web Terminal
# Runs the ttyd daemon
# ==============================================================================
readonly TTYD_COMMAND=(tmux -u new -A -s hassio zsh)
declare -a options

bashio::log.info 'Starting the ttyd daemon...'

# Add default options
options+=(--reconnect 30)

# Set daemon debug level
if bashio::debug; then
    options+=(-d7)
else
    options+=(-d1)
fi

# Listen on socket
options+=(-i /var/run/ttyd.sock)

# Change working directory
cd /root || bashio::exit.nok 'Unable to change working directory'

# Run the ttyd daemon
exec ttyd "${options[@]}" "${TTYD_COMMAND[@]}"
