# ssh/rootfs/etc/services.d/sshd/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Community Hass.io Add-ons: SSH & Web Terminal
# This files check if all user configuration requirements are met
# ==============================================================================
declare -a options

# If SSH is disabled, use a fake sleep process
if bashio::config.false 'ssh.enable'; then
    exec sleep 86400
fi

bashio::log.info 'Starting the SSH daemon...'

# Default options
options+=(-D)
options+=(-e)

# Enable debug mode on the SSH daemon
if bashio::debug; then
  options+=(-d)
fi

# Start the SSH daemon
exec /usr/sbin/sshd "${options[@]}" < /dev/null
