# ssh/rootfs/etc/services.d/sshd/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Home Assistant Community Add-on: SSH & Web Terminal
# Starts the SSH service
# ==============================================================================
declare -a options

# If SSH is disabled, use a fake sleep process
if ! bashio::var.has_value "$(bashio::addon.port 22)"; then
    exec sleep 864000
fi

bashio::log.info 'Starting the SSH daemon...'

# Default options
options+=(-D)
options+=(-e)

# Enable debug mode on the SSH daemon
if bashio::debug; then
  options+=(-d)
fi

# Start the SSH daemon
exec /usr/sbin/sshd "${options[@]}" < /dev/null
