# ssh/rootfs/etc/s6-overlay/s6-rc.d/init-docker/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# shellcheck disable=SC2207
# ==============================================================================
# Home Assistant Community Add-on: Advanced SSH & Web Terminal
# Enables Docker by moving the Docker executable in place.
# ==============================================================================
if bashio::var.false "$(bashio::addon.protected)"; then
    rm /usr/local/bin/docker
    mv /usr/local/bin/.undocked /usr/local/bin/docker
    bashio::log.info 'Docker support has been enabled.'
fi
